# server/deploy/supervisord_conf.ini
[program:Grec]
user=root
command=gunicorn -c deploy/prd_grec.py run:app
autostart=true
autorestart=true
stopasgroup=true
directory=/var/www/Grec/server